# sync-elpa.timer
[Unit]
Description=Runs sync-elpa every 6 hours

[Timer]
OnUnitActiveSec=12h
Unit=sync-elpa.service

[Install]
WantedBy=multi-user.target
